# deploy/herdeck-bridge.service
[Unit]
Description=Herdeck bridge (herdr socket -> Tailscale WebSocket)
After=network-online.target tailscaled.service

[Service]
Environment=HERDR_SOCKET=%h/.config/herdr/herdr.sock
Environment=HERDECK_BIND=100.x.y.z
Environment=HERDECK_PORT=8788
Environment=HERDECK_SERVER_ID=workbox
Environment=HERDECK_TOKEN=change-me
Environment=PYTHONPATH=%h/herdeck/src
ExecStart=%h/herdeck/.venv/bin/python -m herdeck.bridge
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
